# ttrss/service/ttrss/run
#!/bin/bash

exec 2>&1

exec chpst -u www-data /app/update.php --daemon
