# nginx-php/php7/service/php7-fpm/run
#!/bin/bash

exec 2>&1

mkdir -p /run/php

exec /usr/sbin/php-fpm7.0 -c /etc/php/7.0/fpm
